# systemd/xkbset-bouncekeys.service
[Unit]
Description=xkbset bouncekeys

[Service]
ExecStart=%h/.local/bin/xkbset-bouncekeys
Type=oneshot

[Install]
WantedBy=graphical.target
